# sasg_code/scripts/src/HA/sgha-watcher@.service
[Unit]
Description=Watcher for SGHA service

[Service]
Restart=always
WorkingDirectory=/sasg/bin/HA
ExecStart=/sasg/bin/HA/service-watchdog.sh %i
